# sys-libs/jhi/files/jhi
#!/sbin/openrc-run

command=/usr/sbin/jhid
command_background=true
pidfile=/var/run/jhi.pid
name="Intel(R) Dynamic Application Loader Host Interface"

